# etc/systemd/system/lia-firstboot.service
[Unit]
Description=LIA Erststart - Flatpak Installation
After=network-online.target flatpak-system-helper.service
Wants=network-online.target
ConditionPathExists=!/var/lib/lia-firstboot.done

[Service]
Type=oneshot
ExecStart=/usr/libexec/lia-firstboot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
